# workers/gpu/config/supervisor/atos-gpu-worker.conf
[program:atos-gpu-worker]
command=/workspace/atos-gpu-worker/scripts/run_worker.sh
directory=/workspace/atos-gpu-worker
autostart=false
autorestart=true
startsecs=5
stopwaitsecs=30
stdout_logfile=/workspace/logs/worker.log
stderr_logfile=/workspace/logs/worker.err.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
