# workers/gpu/config/supervisor/ollama.conf
[program:ollama]
command=/workspace/scripts/run_ollama.sh
directory=/workspace
autostart=false
autorestart=true
startsecs=5
stopwaitsecs=30
stdout_logfile=/workspace/logs/ollama.log
stderr_logfile=/workspace/logs/ollama.err.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
environment=HOME="/workspace/config/ollama-home",OLLAMA_MODELS="/workspace/ollama-models",OLLAMA_HOST="127.0.0.1:11434"
